# root/etc/services.d/01-yacron/run
#!/usr/bin/with-contenv bash

cd /app || exit

exec \
	s6-setuidgid abc chown -R abc:abc /config && yacron -c /config/cron.yml